# app/server/config/monitor-wifi-watchdog.service
[Unit]
Description=Home Monitor WiFi Watchdog — fallback to setup hotspot if WiFi fails
Documentation=https://github.com/vinu-engineer/rpi-home-monitor
After=network-online.target NetworkManager-wait-online.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no
TimeoutStartSec=120
ExecStart=/opt/monitor/scripts/monitor-wifi-watchdog.sh

[Install]
WantedBy=multi-user.target
