# app/server/config/monitor.service
[Unit]
Description=Home Monitor Server
After=network-online.target nginx.service mediamtx.service local-fs.target
Wants=network-online.target mediamtx.service

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor
ExecStart=/usr/bin/python3 -m flask --app monitor run --host=127.0.0.1 --port=5000
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/monitor
Environment=FLASK_APP=monitor
Environment=MONITOR_DATA_DIR=/data
Environment=MONITOR_RECORDINGS_DIR=/data/recordings
Environment=MONITOR_LIVE_DIR=/data/live
Environment=MONITOR_CONFIG_DIR=/data/config
Environment=MONITOR_CERTS_DIR=/data/certs

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
