# app/camera/config/camera-streamer.service
[Unit]
Description=Camera RTSP Streamer
After=network-online.target avahi-daemon.service local-fs.target NetworkManager.service sys-subsystem-net-devices-wlan0.device
Wants=network-online.target NetworkManager.service

[Service]
Type=simple
User=camera
Group=camera
WorkingDirectory=/opt/camera
# Allow camera user to control the onboard ACT LED for status feedback
ExecStartPre=+/bin/sh -c 'chmod 0666 /sys/class/leds/ACT/trigger /sys/class/leds/ACT/brightness /sys/class/leds/ACT/delay_on /sys/class/leds/ACT/delay_off 2>/dev/null || true'
ExecStart=/usr/bin/python3 -m camera_streamer.main
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/camera
Environment=CAMERA_DATA_DIR=/data
Environment=CAMERA_CONFIG_DIR=/data/config
Environment=CAMERA_CERTS_DIR=/data/certs

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data
PrivateTmp=true
# Camera device access
SupplementaryGroups=video
# Hotspot setup needs NetworkManager D-Bus access and port 80
# CAP_SYS_ADMIN needed to set hostname (mDNS) on read-only rootfs
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
