# app/camera/config/camera-hotspot.service
[Unit]
Description=Camera WiFi Setup Hotspot
Documentation=https://github.com/vinu-engineer/rpi-home-monitor
After=NetworkManager.service local-fs.target sys-subsystem-net-devices-wlan0.device gpio-trigger.service
Wants=NetworkManager.service
Before=camera-streamer.service

# Only start if initial setup has not been completed
ConditionPathExists=!/data/.setup-done

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=90
ExecStart=/opt/camera/scripts/camera-hotspot.sh start
ExecStop=/opt/camera/scripts/camera-hotspot.sh stop

[Install]
WantedBy=multi-user.target
